# docker/loader/crontab
0 0 * * * /usr/local/bin/validol_scraper >> /var/log/validol/scraper.log 2>&1
