# docker/loader/crontab
* * * * * /usr/local/bin/validol_scraper >> /var/log/validol_scraper.log 2>&1
